# rabbitmq/supervisord.conf
[supervisord]
nodaemon=true

[program:rabbit-mq]
command=service rabbitmq-server restart
