# mongodb/supervisord.conf
[supervisord]
nodaemon=true

#[program:env]
#command=/etc/entrypoint.sh

[program:mongodb]
command=/etc/init.d/mongod start
