# elasticsearch/supervisord.conf
[supervisord]
nodaemon=true

[program:elasticsearch]
command=service elasticsearch start
