# redis/supervisord.conf
[supervisord]
nodaemon=true

[program:redis]
command=/etc/init.d/redis start
